# crontab
0 0 * * * ((echo 'Command started' && /usr/local/bin/python3 /usr/src/app/sbin/db_updater.py -v -c -i -v) > /proc/1/fd/1 2>/proc/1/fd/2)
